# dek-rollback.service
[Unit]
Description=Pollek DEK Core Supervisor Rollback
DefaultDependencies=no
After=network.target

[Service]
Type=oneshot
# In a real environment, this might be a small script or a dekctl command
# that restores the .bak file and restarts the service.
ExecStart=/bin/sh -c 'if [ -f /usr/local/bin/dek-core.bak ]; then mv /usr/local/bin/dek-core.bak /usr/local/bin/dek-core; fi; systemctl restart pollek-dek.service'
