# pollek-dek.service
[Unit]
Description=Pollek DEK Core Supervisor
After=network.target

[Service]
Type=notify
ExecStart=/usr/local/bin/dek-core
Restart=on-failure
RestartSec=2

# Health gating and Watchdog
WatchdogSec=30

# Prevent crash-looping (if it crashes 5 times in 60 seconds, stop restarting)
StartLimitIntervalSec=60
StartLimitBurst=5

# Emergency rollback when hitting the StartLimit (unrecoverable new binary)
OnFailure=dek-rollback.service

# Least-privilege capabilities for eBPF
AmbientCapabilities=CAP_BPF CAP_NET_ADMIN

ConfigurationDirectory=pollek-dek
StateDirectory=pollek-dek
LogsDirectory=pollek-dek
RuntimeDirectory=pollek-dek

[Install]
WantedBy=multi-user.target
