# packaging/linux/scripts/postinst
#!/bin/sh
set -e
systemctl daemon-reload || true
systemctl enable pollek-dek || true
systemctl restart pollek-dek || true
