# packaging/linux/scripts/prerm
#!/bin/sh
set -e
systemctl stop pollek-dek || true
systemctl disable pollek-dek || true
