# packaging/linux/pollek-dek.service
[Unit]
Description=Pollek Distributed Enforcement Kernel Core
After=network.target

[Service]
ExecStart=/usr/local/bin/dek-core
Restart=always
RestartSec=5
User=root
Environment=DEK_BOOTSTRAP_PATH=/opt/pollek/bootstrap.json
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
